# roles/rtorrent/templates/rtorrent-supervisor.conf.j2
[program:rtorrent]
command=/usr/local/bin/rtorrent
redirect_stderr=true
user=root
