# roles/rtorrent/templates/rtorrent-supervisor.conf.j2
[program:rtorrent]
command=/usr/local/bin/rtorrent
redirect_stderr=true
environment=HOME=/opt/torrent-files,USER=www-data
directory=/opt/torrent-files
user=www-data
